# kuali-ec2/src/main/resources/jenkins/master/etc/init.d/jenkins
#!/bin/bash
# /etc/init.d/jenkins
# debian-compatible jenkins startup script.
# Jeff Caddel <jcaddel@kuali.org>
#
### BEGIN INIT INFO
# Provides:          jenkins
# Required-Start:    $remote_fs $syslog $network
# Required-Stop:     $remote_fs $syslog $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start jenkins at boot time
# Description:       Controls the jenkins continuous integration engine.
### END INIT INFO

MAXFILES=8192

echo "Service call made - $(date)" > /tmp/jenkins.txt

case "$1" in
  start)
    echo "Starting Jenkins - $(date)" >> /tmp/jenkins.txt
    ulimit -n $MAXFILES >> /tmp/jenkins.txt
    su - jenkins --command '/usr/local/tomcat/bin/tomcat start' >> /tmp/jenkins.txt 2>&1
    ;;
  stop)
    echo "Stopping Jenkins - $(date)" >> /tmp/jenkins.txt
    su - jenkins --command '/usr/local/tomcat/bin/tomcat stop' >> /tmp/jenkins.txt 2>&1
    ;;
  restart)
    echo "Stopping Jenkins - $(date)" >> /tmp/jenkins.txt
    su - jenkins --command '/usr/local/tomcat/bin/tomcat stop' >> /tmp/jenkins.txt 2>&1
    echo "Jenkins stopped - $(date)" >> /tmp/jenkins.txt
    ulimit -n $MAXFILES >> /tmp/jenkins.txt
    echo "Starting Jenkins - $(date)" >> /tmp/jenkins.txt
    su - jenkins --command '/usr/local/tomcat/bin/tomcat start' >> /tmp/jenkins.txt 2>&1
    echo "Jenkins started - $(date)" >> /tmp/jenkins.txt
    ;;
  *)
    echo "Usage: jenkins {start|stop|restart}" >&2
    exit 3
    ;;
esac

exit 0
